# old backup/systemd/security-agent.service
# systemd/security-agent.service
[Unit]
Description=AI Security Agent
After=network.target mcp-server.service ollama.service
Requires=mcp-server.service

[Service]
Type=simple
User=your-user
WorkingDirectory=/path/to/security-automation
Environment="PATH=/path/to/security-automation/venv/bin"
ExecStart=/path/to/security-automation/venv/bin/python3 agent/security_agent.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target